# services.d/docker-gen/run
#!/bin/bash
/usr/local/bin/docker-gen -watch -notify -endpoint unix:///tmp/docker.sock "nginx -s reload" /app/nginx.tmpl /etc/nginx/conf.d/default.conf